# continuonbrain/systemd/continuonbrain-startup.service
# ContinuonBrain Startup Service
#
# Installation:
#   1. Copy this file to user systemd directory:
#      cp continuonbrain-startup.service ~/.config/systemd/user/
#
#   2. Enable auto-start on boot (without requiring user login):
#      sudo loginctl enable-linger $USER
#
#   3. Enable and start the service:
#      systemctl --user enable continuonbrain-startup.service
#      systemctl --user start continuonbrain-startup.service
#
#   4. Check status:
#      systemctl --user status continuonbrain-startup.service
#
# Note: Adjust paths in Environment variables below for your deployment.

[Unit]
Description=ContinuonBrain Robot startup (health checks + Robot API server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust these paths for your deployment before installing:
Environment=PYTHONPATH=/opt/continuonos/ContinuonXR
Environment=CONFIG_DIR=/opt/continuonos/brain
# Prefer headless mode and skip background trainer on Pi5 by default
Environment=CONTINUON_HEADLESS=1
Environment=CONTINUON_ENABLE_BACKGROUND_TRAINER=0
WorkingDirectory=/opt/continuonos/ContinuonXR
User=pi
ExecStart=/opt/continuonos/ContinuonXR/.venv/bin/python3 -m continuonbrain.startup_manager --config-dir ${CONFIG_DIR}
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
