# continuonbrain/systemd/continuonbrain-startup.service
# Continuon Brain runtime autostart (systemd)
#
# Goal: start the Continuon Brain runtime on boot and automatically restart it if it crashes.
# This is a *system* unit (installed under /etc/systemd/system), not a per-login user unit.
#
# Install:
#   sudo mkdir -p /etc/continuonbrain
#   sudo cp continuonbrain/systemd/continuonbrain.env.example /etc/continuonbrain/continuonbrain.env
#   sudo cp continuonbrain/systemd/continuonbrain-startup.service /etc/systemd/system/continuonbrain-startup.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now continuonbrain-startup.service
#
# Logs:
#   sudo journalctl -u continuonbrain-startup.service -f
#
# Notes:
# - "Wake" behavior: the service typically continues through suspend; if it stops, Restart=always brings it back.
# - Prefer editing /etc/continuonbrain/continuonbrain.env rather than modifying this unit.

[Unit]
Description=Continuon Brain runtime startup (health checks + Robot API server)
After=network-online.target
Wants=network-online.target
# Disable start burst limiting for a “keep trying” appliance behavior.
StartLimitIntervalSec=0

[Service]
Type=forking
User=root
Group=root

# Site overrides (recommended).
EnvironmentFile=-/etc/continuonbrain/continuonbrain.env

# Defaults (can be overridden in /etc/continuonbrain/continuonbrain.env)
Environment=CONTINUON_REPO=/opt/continuonos/ContinuonXR
Environment=CONFIG_DIR=/opt/continuonos/brain
Environment=PYTHONPATH=/opt/continuonos/ContinuonXR
Environment=CONTINUON_PYTHON=/usr/bin/python3
Environment=PIDFILE=/opt/continuonos/brain/runtime_pids/startup_manager.pid

ExecStartPre=/usr/bin/test -d $CONTINUON_REPO
ExecStart=/usr/bin/env bash -lc '${CONTINUON_REPO}/scripts/start_rpi.sh start'
ExecStop=/usr/bin/env bash -lc '${CONTINUON_REPO}/scripts/start_rpi.sh stop'
PIDFile=/opt/continuonos/brain/runtime_pids/startup_manager.pid

# Make "restart" robust (covers crashes and stop/start cycles).
Restart=always
RestartSec=2

# Graceful shutdown.
KillSignal=SIGINT
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
